# etc/cron.d/apc
*/5 * * * * root /usr/bin/curl -so /tmp/apc.php http://oscar.rewse.jp/admin/apc.php
